# images-build/common/rc.local
#!/bin/sh
mkdir -p /root/.ssh
cp /vagrant/secret/id_rsa.pub /root/.ssh/authorized_keys
chown -R root:root /root/.ssh
chmod 0700 /root/.ssh
chmod 0600 /root/.ssh/authorized_keys
true
